# dove6-package/rootfs/etc/systemd/system/dove6.service
[Unit]
Description=Dove6 Passenger Information Display
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=aviarail
WorkingDirectory=/opt/dove6
ExecStart=/opt/dove6/bin/start-dove6
Restart=always
RestartSec=3

NoNewPrivileges=true
PrivateTmp=true
ProtectHome=false

[Install]
WantedBy=multi-user.target
